# dht11serv.service
[Unit]
Description=Rust Application Service
After=network.target

[Service]
ExecStart=/usr/share/dht11serv/dht11serv/target/release/dht11serv
Restart=always

[Install]
WantedBy=multi-user.target